# pacUpdt.timer
[Unit]
Description=Run PacUpdt every 30 min

[Timer]
OnBootSec=60sec
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
