# pacUpdt.service
[Unit]
Description=Writes number of available updates to text file

[Service]
Type=simple
ExecStart=/usr/bin/pacUpdt
Restart=no
